# files/conserver.service
[Unit]
Description=Conserver Serial-Port Console Daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/conserver -d -i -u -M 127.0.0.1

[Install]
WantedBy=multi-user.target
